# crontab
* * * * * cd /home/che0520/ubike && /usr/bin/python3 webToMysql.py >> run.log